# dist/src/main/unix/scripts/postinst
#!/bin/sh

set -e

#DEB

GROUP='stratio'
USER='stratio'
DIR_INSTALL='/opt/sds/notebook'
DIR_CONF='/etc/sds/notebook'
DIR_LOG='/var/log/sds/notebook'
DIR_RUN='/var/run/sds'


chmod 755 /etc/init.d/explorer
chmod 755 /opt/sds/explorer/bin
chmod +x /opt/sds/explorer/bin/*
chmod +x /etc/default/explorer-functions
chmod +x /etc/default/explorer-variables

case "$1" in
  configure)

  # Add stratio group if it doesn't already exists
  if ! getent group $GROUP >/dev/null 2>&1; then
    groupadd ${GROUP}
  fi

  # Add stratio user if it doesn't already exists
  if ! getent passwd $USER >/dev/null 2>&1; then
    [ -d /opt ] || mkdir /opt
    useradd -M -d /opt/sds -g ${GROUP} -s /bin/bash ${USER}
  else
    usermod -G $GROUP $USER
  fi

  if [ ! -d "$DIR_LOG" ]; then
    mkdir $DIR_LOG
  fi

  chown root:root /etc/init.d/explorer

  chown -R ${USER}:${GROUP} $DIR_INSTALL
  chown -R ${USER}:${GROUP} $DIR_CONF
  chown -R ${USER}:${GROUP} $DIR_LOG
  chmod -R 775 $DIR_INSTALL
  chown ${USER}:${GROUP} $DIR_RUN
  chmod 775 $DIR_LOG

  update-rc.d explorer defaults
  ;;
esac


exit 0